# gate/etc/systemd/system/wg-quick@.service.d/override.conf
[Unit]
After=nftables.service docker.service
Requires=nftables.service docker.service
