# spire/usr/lib/systemd/system/hass-bridge.socket
[Socket]
ListenStream=/var/lib/hass-bridge/socket
SocketUser=mendess
SocketGroup=mendess
SocketMode=0660

[Install]
WantedBy=sockets.target
