# spire/usr/lib/systemd/system/start-gaming-mode.service
[Service]
Type=oneshot
User=mendess
ExecStart=/usr/local/bin/start-gaming-mode
